# elasticsearch/curator.cron
@daily root curator delete --disk-space=8 indices
